# deploy/modelscope/supervisord.conf
; Another Me — ModelScope 创空间进程编排
; 单容器内顺序拉起：postgres → sandbox-runner → backend → nginx。
; 所有日志统一打到 stdout/stderr，便于在创空间「查看日志」中查看。
[supervisord]
nodaemon=true
user=root
logfile=/dev/null
logfile_maxbytes=0
pidfile=/run/supervisord.pid

; 内置 Postgres：首次启动初始化数据目录到 /mnt/workspace/pgdata（见脚本）
[program:postgres]
command=/usr/local/bin/run-postgres.sh
autostart=true
autorestart=true
startretries=10
priority=10
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

; sandbox-runner：复用 backend 的 venv（依赖完全一致），仅监听本地回环
[program:sandbox]
command=uvicorn main:app --host 127.0.0.1 --port 8001
directory=/app/sandbox
autostart=true
autorestart=true
priority=20
environment=PATH="/app/backend/.venv/bin:%(ENV_PATH)s"
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

; backend：等待 Postgres 就绪 → alembic upgrade → seed → uvicorn :8000
[program:backend]
command=/usr/local/bin/run-backend.sh
autostart=true
autorestart=true
priority=30
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

; nginx：对外唯一入口 :7860（先起无妨，backend 未就绪时 /api 暂时 502）
[program:nginx]
command=nginx -g 'daemon off;'
autostart=true
autorestart=true
priority=40
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
